# dronecode-camera-manager.service.in
[Unit]
Description=Dronecode Camera Manager
After=multi-user.target

[Service]
Type=simple
ExecStart=@bindir@/dcm
Restart=on-failure

[Install]
WantedBy=multi-user.target
